# docker/build/assets/supervisord.conf
[supervisord]
nodaemon=true

[program:postfix]
command=/usr/local/bin/iwfpostfix.sh

[program:phpfpm]
command=/usr/local/sbin/php-fpm
