# docker/build/assets/etc/supervisor/conf.d/iwfsupervisor.conf
[supervisord]
logfile=/dev/stderr                             ; supervisord log file
loglevel=warn                                   ; info, debug, warn, trace
pidfile=/var/run/supervisord.pid                ; pidfile location
nodaemon=true                                   ; run supervisord as a daemon
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
user=root                                       ; default user
childlogdir=/var/log/supervisor/               ; where child log files will live
logfile_maxbytes=0

[program:init]
command=/bin/bash -c "/usr/bin/sudo SUDO_FORCE_REMOVE=yes /usr/bin/apt-get -y remove sudo netcat-openbsd | sed 's/^/[init] /'"
stderr_logfile = /dev/stderr
stdout_logfile = /dev/stdout
autorestart=false
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0

[program:phpfpm]
command=/bin/bash -c "/usr/local/sbin/php-fpm | sed 's/^/[phpfpm] /'"
stderr_logfile = /dev/stderr
stdout_logfile = /dev/stdout
autorestart=unexpected
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0

[program:cron]
command=/bin/bash -c "/usr/sbin/cron -f -L 7 | sed 's/^/[cron] /'"
stderr_logfile = /dev/stderr
stdout_logfile = /dev/stdout
autorestart=unexpected
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0
