# docker/build/assets/etc/supervisor/conf.d/iwfsupervisor.conf
[supervisord]
logfile=/dev/stdout                             ; supervisord log file
loglevel=warn                                   ; info, debug, warn, trace
pidfile=/var/run/supervisord.pid                ; pidfile location
nodaemon=true                                   ; run supervisord as a daemon
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
user=root                                       ; default user
logfile_maxbytes=0

[program:init]
command=/bin/bash -c "/usr/bin/sudo SUDO_FORCE_REMOVE=yes /usr/bin/apt-get -y remove sudo netcat-openbsd 2>&1 | sed --unbuffered 's/^/[init] /' && sleep 60"
stderr_logfile = /dev/stderr
stdout_logfile = /dev/stdout
autorestart=false
startretries=0
exitcodes=0
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0

[program:phpfpm]
command=/bin/bash -c "/usr/local/sbin/php-fpm 2>&1 | sed --unbuffered 's/^/[phpfpm] /'"
stderr_logfile = /dev/stderr
stdout_logfile = /dev/stdout
autorestart=unexpected
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0

[program:cron]
command=/bin/bash -c "/usr/sbin/cron -f -L 7 2>&1 | sed --unbuffered 's/^/[cron] /'"
stderr_logfile = /dev/stderr
stdout_logfile = /dev/stdout
autorestart=unexpected
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0
